# launchd/com.dotfiles.homebrew_monthly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
    Файл com.dotfiles.homebrew_monthly.plist настраивает запуск скрипта
    homebrew_monthly_update.sh раз в месяц.
    Скрипт выполняет обновление Homebrew (formulae & casks) и очистку кэша.

    Инструкция: создайте символическую ссылку этого файла в каталоге
    ~/Library/LaunchAgents для активации:
        ln -s ~/dotfiles/launchd/com.dotfiles.homebrew_monthly.plist \
              ~/Library/LaunchAgents/com.dotfiles.homebrew_monthly.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dotfiles.homebrew_monthly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/khurshidddbek/dotfiles/scripts/homebrew_monthly_update.sh</string>
    </array>

    <!-- Запуск 1-го числа каждого месяца в 15:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Логируем stdout и stderr в общий файл -->
    <key>StandardOutPath</key>
    <string>/Users/khurshidddbek/dotfiles/biweekly_update.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/khurshidddbek/dotfiles/biweekly_update.log</string>

    <!-- Резервный ежемесячный интервал (30 дней), если Mac был выключен -->
    <key>StartInterval</key>
    <integer>2592000</integer>
</dict>
</plist> 