# launchd/com.dotfiles.monthly_tasks.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
    Файл com.dotfiles.monthly_tasks.plist ежедневно в 12:00 запускает
    anacron-обёртку run_monthly_tasks_anacron.sh. Обёртка сама решает,
    пора ли выполнять ежемесячные задачи (run_monthly_tasks.sh): они
    запускаются, если с последнего успешного прогона прошло 30 дней.
    Такая схема не теряет запуск, если мак был выключен в назначенный
    момент — задачи выполнятся в ближайший включённый день.
-->
<plist version="1.0">
<dict>
    <!-- Уникальный идентификатор задачи -->
    <key>Label</key>
    <string>com.dotfiles.monthly</string>

    <!-- Команда для выполнения -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/khurshidddbek/dotfiles/scripts/scheduled/run_monthly_tasks_anacron.sh</string>
    </array>

    <!-- Ежедневная проверка в 12:00 (сам прогон — раз в 30 дней) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>12</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Логи -->
    <key>StandardOutPath</key>
    <string>/Users/khurshidddbek/dotfiles/scheduled_tasks.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/khurshidddbek/dotfiles/scheduled_tasks.log</string>
</dict>
</plist>
